# door-sensor.service
[Unit]
Description=Door Sensor Daemon
After=network.target
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/hitl-tools/src/
ExecStart=/usr/bin/python3 /opt/hitl-tools/src/door_sensor_daemon.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
AmbientCapabilities=CAP_SYS_RAWIO
ProtectHome=true
ProtectSystem=strict
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
